# packaging/dux.service
[Unit]
Description=dux realtime disk-usage + file-search index
Documentation=man:dux(1)
After=local-fs.target

[Service]
Type=simple
# Build the index ONLY if one doesn't exist yet (first install / fresh boot).
# On restarts the daemon already maintains a live index, so we skip the scan —
# that keeps `systemctl restart dux` instant instead of blocking ~30-60s.
# Force a fresh full scan any time with: dux scan /   (or systemctl start dux-scan).
ExecStartPre=/bin/sh -c '[ -s /var/lib/dux/dux.db ] || /usr/bin/dux scan / --quiet'
ExecStart=/usr/bin/dux daemon /
Restart=on-failure
RestartSec=5
# Liveness heartbeat lives in /run/dux (tmpfs), not SQLite, so an idle daemon
# makes zero DB/WAL writes. RuntimeDirectory creates it and removes it on stop,
# so the heartbeat file vanishes when the daemon dies (no stale "live").
RuntimeDirectory=dux
RuntimeDirectoryMode=0755
# First-time scan can take a while on big/slow disks; give it room.
TimeoutStartSec=3600

# Stay out of the way of real workloads.
Nice=10
IOSchedulingClass=idle
CPUWeight=20
# The scan buffers all node metadata in memory; budget generously and prefer a
# soft limit so very large filesystems don't OOM mid-scan.
MemoryHigh=2G
MemoryMax=4G

# fanotify needs CAP_SYS_ADMIN; reading all files needs DAC override.
# NOTE: do NOT enable ProtectSystem/ProtectHome/PrivateTmp etc. — those put the
# service in a private mount namespace, and a fanotify filesystem mark there
# would not see writes happening through the host's mounts.
AmbientCapabilities=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH

[Install]
WantedBy=multi-user.target
