# packaging/dux.service
[Unit]
Description=dux realtime disk-usage + file-search index
Documentation=man:dux(1)
After=local-fs.target

[Service]
Type=simple
# Build the index ONLY if one doesn't exist yet (first install / fresh boot).
# On restarts the daemon already maintains a live index, so we skip the scan —
# that keeps `systemctl restart dux` instant instead of blocking ~30-60s.
# Force a fresh full scan any time with: dux scan /   (or systemctl start dux-scan).
ExecStartPre=/bin/sh -c '[ -s /var/lib/dux/dux.db ] || /usr/bin/dux scan / --quiet'
ExecStart=/usr/bin/dux daemon /
Restart=on-failure
RestartSec=5
# Liveness heartbeat lives in /run/dux (tmpfs), not SQLite, so an idle daemon
# makes zero DB/WAL writes. RuntimeDirectory creates it and removes it on stop,
# so the heartbeat file vanishes when the daemon dies (no stale "live").
RuntimeDirectory=dux
RuntimeDirectoryMode=0755
# First-time scan can take a while on big/slow disks; give it room.
TimeoutStartSec=3600

# Stay out of the way of real workloads.
Nice=10
IOSchedulingClass=idle
CPUWeight=20
# The scan buffers all node metadata in memory; budget generously and prefer a
# soft limit so very large filesystems don't OOM mid-scan.
MemoryHigh=2G
MemoryMax=4G

# fanotify needs CAP_SYS_ADMIN; reading all files needs DAC override.
# NOTE: do NOT enable ProtectSystem/ProtectHome/PrivateTmp/PrivateMounts/
# ReadOnlyPaths etc. — those put the service in a PRIVATE MOUNT NAMESPACE, and a
# fanotify filesystem mark there would not see writes through the host's mounts.
AmbientCapabilities=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
# Cap what the process can ever acquire to exactly the two it needs (CAP_SYS_ADMIN
# is broad, so this is the main lever for reducing the privilege surface).
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
# Hardening that does NOT require a private mount namespace (all fanotify-safe):
NoNewPrivileges=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
# The alert exec (--alert-exec) and nothing else needs the network; if you don't
# use it you can also set: RestrictAddressFamilies=AF_UNIX

# --- index exposure (READ THIS) -------------------------------------------
# The index at /var/lib/dux/dux.db contains the FULL FILENAME LIST of the root
# filesystem, including names inside directories a user could not otherwise
# traverse. By default it is world-readable (0644) so unprivileged `dux find` /
# `dux status` work (dux is a companion to user-facing tooling). This is an
# information-disclosure tradeoff. To make the index ROOT-ONLY instead, set:
#     StateDirectory=dux
#     StateDirectoryMode=0700
#     UMask=0077
# (the daemon then creates dux.db 0600); unprivileged users lose direct query
# access and must go through a privileged wrapper.

[Install]
WantedBy=multi-user.target
